# crontab
# m h dom mon dow user  command

0 * * * * root /bin/bash /backup.sh
